# openthread_border_router/rootfs/etc/s6-overlay/s6-rc.d/dhcpd-daemon/finish
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
#==============================================================================
# OpenThread BorderRouter Daemon finish script
#==============================================================================


/var/lib//usr/sbin/dhcpcd -X